# shell/amount/55-sd.rules
# Mount serial disks
KERNEL=="sd[a-z][0-9]*", RUN+="/usr/local/share/udev/sd.action"
