# systemd/echo/echo.socket
[Unit]
Description=Echo Socket

[Socket]
ListenStream=1234
Accept=yes

[Install]
WantedBy=sockets.target
